# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/api/artisan schedule:run >> /var/log/cron.log 2>&1

* * * * * laradock /bin/sh -c 'cd /var/www/api/ && git add . && git stash && git pull --all && /usr/local/bin/composer install' >> /var/log/cron.log 2>&1

* * * * * root supervisorctl restart all >> /var/log/cron.log 2>&1

* * * * * root /bin/sh -c 'chown -R laradock /var/www/api/' >> /var/log/cron.log 2>&1

0 12 * * * laradock /bin/sh -c '> /var/log/cron.log' >> /dev/null 2>&1
